# node-ruby/init_script
### BEGIN INIT INFO
# Provides:          blusteromatic_node
# Required-Start:    $local_fs $remote_fs
# Required-Stop:     $local_fs $remote_fs
# Default-Start:     2 3 4 5
# Default-Stop:      S 0 1 6
# Short-Description: blusteromatic node initscript
# Description:       blusteromatic
### END INIT INFO

# Do NOT "set -e"

DAEMON="cd blusteromatic/node-ruby && ./node.rb"
USER="blusteromatic"
#SCRIPT_NAME="/etc/init.d/foreverb-username"
SCRIPT_NAME="$0"

case "$1" in
  start)
  su -l $USER -c "$DAEMON start"
  ;;
  stop)
  su -l $USER -c "$DAEMON stop"
  ;;
  restart)
  su -l $USER -c "$DAEMON restart"
  ;;
  status)
  su -l $USER -c "foreverb list"
  ;;
  log)
  su -l $USER -c "foreverb tail node.rb"
  ;;
  *)
  echo "Usage: $SCRIPT_NAME {start|stop|restart}" >&2
  exit 3
  ;;
esac
